# homi.service.in
[Unit]
Description=Host Orchestrated Multipath I/O

[Service]
Type=simple
RuntimeDirectory=homi
StandardOutput=journal
StandardError=journal
ExecStart=@bindir@/homid --config @config@

[Install]
WantedBy=multi-user.target
